# systemd/agnetic-staragent.service
[Unit]
Description=Agnetic Starship OS - StarAgent Telemetry Collector
After=network.target agnetic-nats.service
Requires=agnetic-nats.service
Documentation=https://github.com/andromi-hash/agnetic-os

[Service]
Type=simple
ExecStart=/opt/agnetic/bin/staragent
Restart=always
RestartSec=5
User=agnetic
Group=agnetic
Environment=RUST_LOG=info
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/agnetic /tmp
ProtectHome=true
KillMode=mixed
TimeoutStopSec=15
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
