# systemd/agnetic-agents.service
[Unit]
Description=Agnetic OS Agent Daemons (Proxy, Romi, Ergo)
Documentation=https://github.com/andromi-hash/agnetic-os
After=agnetic-nats.service
BindsTo=agnetic-nats.service

[Service]
Type=simple
ExecStart=/home/tech/.hermes/hermes-agent/venv/bin/python3 /home/tech/agnetic-os/agents/agent_daemon.py
Restart=always
RestartSec=5
User=tech
Group=tech
Environment=NATS_URL=nats://127.0.0.1:4222

[Install]
WantedBy=multi-user.target
